# scripts/spv_vnc.conf
# From https://github.com/ct2034/docker-ubuntu-vnc-desktop
[supervisord]
redirect_stderr=true
stopsignal=QUIT
autorestart=true
directory=/root

[group:x]
programs=xvfb,wm,x11vnc

[program:xvfb]
priority=10
command=/usr/bin/Xvfb :1 -screen 0 1440x900x24
stopsignal=KILL

[program:wm]
priority=15
command=/usr/bin/openbox
environment=DISPLAY=":1",HOME="/root",USER="root"

[program:x11vnc]
priority=20
command=x11vnc -display :1 -xkb -forever -shared -repeat -rfbport 5901
